# docker/wis2-gdc-management.cron
0 0 * * * su -c "wis2-gdc archive /data/wis2-discovery-metadata-archive.zip" > /proc/1/fd/1 2>/proc/1/fd/2
