# wis2-gdc-management/docker/wis2-gdc-management.cron
0 0 * * * /venv/bin/wis2-gdc archive $WIS2_GDC_METADATA_ARCHIVE_ZIPFILE
30 * * * * /venv/bin/pywcmp bundle sync
