# wis2-gdc-management/docker/wis2-gdc-management.cron
0 0 * * * wis2-gdc archive $WIS2_GDC_METADATA_ARCHIVE_ZIPFILE
